# eperp2api-django-q-scheduler.timer
[Unit]
Description=Timer for Django Q Scheduler service

[Timer]
OnBootSec=30s

[Install]
WantedBy=timers.target
