# gunicorn-eperp2api.timer
[Unit]
Description=Timer for gunicorn service

[Timer]
OnBootSec=30s

[Install]
WantedBy=timers.target
